# expo-updates-server/deploy/supervisor/updates-server.conf
; ── Custom OTA Updates Server ──────────────────────────────────────────────────
; Copy this file to /etc/supervisor/conf.d/updates-server.conf
; Replace all <PLACEHOLDER> values with your actual configuration.
;
; Values to replace:
;   <USER>            - system user (e.g. otaupdates)
;   <PROJECT_DIR>     - absolute path to expo-updates-server/
;   <HOSTNAME_URL>    - public URL (e.g. https://updates.example.com:3000)
;   <PUBLISH_API_KEY> - secret key for POST /api/publish
;   <OTA_CLEANUP_ENABLED> - true to run detached cleanup after successful publish
;   <OTA_CLEANUP_KEEP_LATEST> - number of updates to retain per runtimeVersion
;   <SSL_CERT_PATH>   - path to SSL certificate
;   <SSL_KEY_PATH>    - path to SSL private key
;   <SSL_CA_PATH>     - path to SSL CA bundle
;
[program:updates-server]
command=/usr/bin/env node server.js
directory=<PROJECT_DIR>
user=<USER>
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true

environment=
    NODE_ENV="production",
    PORT="3000",
    HOSTNAME="<HOSTNAME_URL>",
    PUBLISH_API_KEY="<PUBLISH_API_KEY>",
    OTA_CLEANUP_ENABLED="<OTA_CLEANUP_ENABLED>",
    OTA_CLEANUP_KEEP_LATEST="<OTA_CLEANUP_KEEP_LATEST>",
    SSL_CERT_PATH="<SSL_CERT_PATH>",
    SSL_KEY_PATH="<SSL_KEY_PATH>",
    SSL_CA_PATH="<SSL_CA_PATH>"

stdout_logfile=<PROJECT_DIR>/logs/supervisor/updates-server.out.log
stderr_logfile=<PROJECT_DIR>/logs/supervisor/updates-server.err.log
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile_backups=5
